# systemd/yt-email-ingestion.service
[Unit]
Description=YouTube/Media Email Ingestion
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=YOUR_USERNAME
Group=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/projects/yt-email-ingestion
ExecStart=/home/YOUR_USERNAME/projects/yt-email-ingestion/.venv/bin/python /home/YOUR_USERNAME/projects/yt-email-ingestion/scripts/mail_fetch.py

[Install]
WantedBy=multi-user.target
