# systemd/yt-email-ingestion.service
[Unit]
Description=YouTube/Media Email Ingestion
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=nick
Group=nick
WorkingDirectory=/home/nick/projects/yt-email-ingestion
ExecStart=/home/nick/projects/yt-email-ingestion/.venv/bin/python /home/nick/projects/yt-email-ingestion/scripts/mail_fetch.py

[Install]
WantedBy=multi-user.target
